# piban-broker/dist/piban.service
[Unit]
Description=My Piban broker Service
After=multi-user.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi
Environment="MQTTHOST=127.0.0.1"
ExecStart=/home/pi/piban.py

[Install]
WantedBy=multi-user.target